# dns-proxy-supervisor.conf
[program:dns-proxy]
command=python3 -u /usr/bin/dns-proxy
user=root
redirect_stderr=true
stdout_logfile=/dev/stdout
startsecs=10
autostart=true
autoresart=true
stdout_logfile_maxbytes=0
